# scripts/server.service
[Unit]
Description= hugo server

[Service]

[Install]